# misc/systemd/system/deepin-update-manager.service
[Unit]
Description=deepin Update Manager

[Service]
Type=simple
User=deepin-update-manager
ExecStart=/usr/libexec/deepin-update-manager
